# timescaledb/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/usr/bin/with-contenv bashio
exec chpst -u postgres postgres -D /data/postgres -k /tmp
